# debian/cron.daily
#!/bin/sh
set -e
if [ -e /etc/etckeeper/daily ]; then
	/etc/etckeeper/daily
fi
